# jackman_plug@.service
[Unit]
Description=automatically switches jackdbus master device
After=suspend.target

[Service]
Type=forking
ExecStart=/usr/bin/jackman_udev_plug %I &
Restart=on-failure
RestartSec=5

[Install]
WantedBy=suspend.target
WantedBy=default.target
